# templates/seaweedfs-master.service.erb
[Unit]
Description=SeaweedFS master
Requires=network-online.target
After=network-online.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Restart=on-failure
ExecStart=<%= @bin_dir %>/<%= @bin_name %> -logdir=<%= @log_dir %> master \
  -mdir=<%= @metadata_dir %> \
  -port=<%= @port %> \
  -ip=<%= @ipaddress %> \
  -defaultReplication=<%= @replication %> \
  -volumeSizeLimitMB=<%= @volume_size_limit %> \
  -volumePreallocate
  <%- unless @peers.empty? -%> -peers=<%= @peers.join(',') -%> \<%- end -%> 
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
