# templates/seaweedfs-filer.service.erb
[Unit]
Description=SeaweedFS filer
Requires=network-online.target
After=network-online.target

[Service]
Restart=on-failure
ExecStart=<%= @bin_dir %>/<%= @bin_name %> -logdir=<%= @log_dir %> filer -ip=<%= @ipaddress %> -port=<%= @port %> -master=<%= @master_ip %>:<%= @master_port %> -defaultReplicaPlacement=<%= @replication %>
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
